# build/overlay/etc/systemd/system/rsync-zram.timer
[Unit]
Description=Periodic ZRAM to Disk Sync
Requires=zram-log.service zram-data.service

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Persistent=true

[Install]
WantedBy=timers.target
